# vendor/gitlab.connectwisedev.com/platform/platform-tasking-service/services/continuum/config/platform-tasking-service.conf

#Continuum Tasking Service Microservice
#description Tasking Service Microservice

start on runlevel [2345]
stop on runlevel [06]

limit nofile 30000 30000

env cfg_file=/opt/continuum/config/ctm_tasking_cfg.json
env daemon=/opt/continuum/taskingService/platform-tasking-service
env logDir=/opt/continuum/log

#Keep service alive upon failure,track forked process, respawn but no more than 10 times in 5 seconds 
expect fork
respawn
respawn limit 10 5

pre-start script
  mkdir -p $logDir
end script

chdir /opt/continuum/taskingService
exec $daemon -config=$cfg_file &
